# systemd/lunahub.service
[Unit]
Description=LunaHub Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lunahub
Group=lunahub
Environment=LUNAHUB_CONFIG=/etc/lunahub/config.json
ExecStart=/usr/local/bin/lunahub serve
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/lunahub /var/log/lunahub

[Install]
WantedBy=multi-user.target
